# templates/program.conf.erb
[program:<%= @name %>]
command=<%= @command %>
<% if @user -%>
user=<%= @user %>
<% end -%>
<% if @directory -%>
directory=<%= @directory %>
<% end -%>
<% if @environment -%>
environment=<%= @environment %>
<% end -%>
